# build_dir/target-mipsel_mips-interAptiv_musl/root.orig-siflower/usr/lib/opkg/info/gl-repeater.list
/usr/bin/gl_health
/etc/hotplug.d/wifi/repeater
/etc/init.d/gl_health
/etc/init.d/gl_ssids
