# build_dir/target-mipsel_mips-interAptiv_musl/root.orig-siflower/usr/lib/opkg/info/rtty-openssl.list
/etc/config/rtty
/etc/uci-defaults/50_rtty
/usr/sbin/rtty
/etc/init.d/rtty
