# build_dir/target-mipsel_mips-interAptiv_musl/root-siflower/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/etc/banner.failsafe
/usr/bin/sfwifi_ate
/etc/resolv.conf
/etc/init.d/done
/bin/config_generate
/var
/lib/preinit/00_preinit.conf
/etc/init.d/system
/etc/preinit
/etc/atetools/power_save_ant2.txt
/etc/localtime
/etc/atetools/sfwifi_ate
/lib/upgrade/platform.sh
/etc/shells
/etc/board.d/01_network
/lib/preinit/70_initramfs_test
/bin/config_check.sh
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/preinit/10_indicate_preinit
/usr/libexec/login.sh
/bin/hnat_set_single_entry.sh
/usr/bin/ate_init_server
/etc/fstab
/usr/lib/os-release
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/shadow
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/etc/board.d/99-default_network
/lib/preinit/10_indicate_failsafe
/usr/bin/hnat_update_interface.sh
/etc/os-release
/bin/sf_reset.sh
/etc/config/bindtable
/etc/atetools/ate_init_server
/lib/preinit/80_mount_root
/lib/preinit/81_urandom_seed
/etc/rc.local
/lib/functions/uci-defaults.sh
/usr/bin/ate_cmd
/etc/init.d/led
/etc/ethers
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/opkg/keys/7a9fb67e92b98f23
/etc/rc.button/failsafe
/lib/upgrade/fwtool.sh
/etc/atetools/power_save.txt
/etc/services
/lib/functions/system.sh
/usr/bin/hnat_read.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/init.d/urandom_seed
/etc/opkg/distfeeds.conf
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/atetools/ate_cmd
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/lib/upgrade/stage2
/etc/uci-defaults/13_fix_group_user
/sbin/firstboot
/usr/bin/ate_init_server.sh
/etc/hotplug.d/net/00-sysctl
/etc/init.d/wan_speed
/etc/hosts
/lib/functions.sh
/etc/atetools/ate_init
/lib/functions/leds.sh
/bin/wanLinkStatus
/etc/atetools/ate_init_server.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/atetools/gain_offset.txt
/etc/profile
/usr/bin/ate_server
/usr/bin/hnat_enable.sh
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/bin/hnat_set_multi_entry.sh
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/sysctl.d/10-default.conf
/etc/passwd
/lib/preinit/02_default_set_state
/bin/test_multi_send.sh
/etc/mtab
/etc/init.d/boot
/etc/config/btn_ctrl_cfg
/usr/bin/ate_init
/etc/config/style
/sbin/urandom_seed
/lib/functions/preinit.sh
/etc/config/siwifi
/bin/ipcalc.sh
/etc/diag.sh
/lib/667.img
/usr/bin/power_save.txt
/etc/atetools/power_save_ant1.txt
