# build_dir/target-mipsel_mips-interAptiv_musl/root-siflower/usr/lib/opkg/info/gl-tertf.list
/etc/config/tertf
/etc/init.d/gl_tertf
/lib/upgrade/keep.d/tertf
/usr/bin/gltertf
