# build_dir/target-mipsel_mips-interAptiv_musl/root.orig-siflower/etc/uci-defaults/fix-network-wan6
uci -q get network.wan6 > /dev/null && exit 0

uci set network.wan6=interface
uci set network.wan6.proto=dhcpv6
uci commit network
