# build_dir/target-mipsel_mips-interAptiv_musl/root-siflower/etc/uci-defaults/50_rtty
[ -f /etc/wbf/customer ] || exit 0

uci set rtty.general.ssh_en=1
uci set rtty.general.web_en=1
uci commit