# build_dir/target-mipsel_mips-interAptiv_musl/root-siflower/etc/uci-defaults/fix-non-admin-user
uci -q get glconfig.user > /dev/null || {
uci batch <<-EOF
	set glconfig.user=login
	commit glconfig
EOF
}
